# systemd/fb-messenger-bridge.service
[Unit]
Description=Hermes <-> Facebook Messenger bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/fb-messenger-bridge
EnvironmentFile=-/opt/fb-messenger-bridge/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/fb-messenger-bridge/venv/bin/python -u /opt/fb-messenger-bridge/bridge.py
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5
StandardOutput=append:/opt/fb-messenger-bridge/service.log
StandardError=append:/opt/fb-messenger-bridge/service.log

[Install]
WantedBy=multi-user.target
